# deploy/systemd/consulting-site.service
[Unit]
Description=Consulting Site Gunicorn Service
After=network.target

[Service]
User=deploy
Group=www-data
WorkingDirectory=/opt/consulting-site/current
EnvironmentFile=-/opt/consulting-site/shared/env/app.env
UMask=0027

ExecStart=/opt/consulting-site/current/.venv/bin/gunicorn \
  --workers 3 \
  --bind 127.0.0.1:8000 \
  --access-logfile - \
  --error-logfile - \
  wsgi:app

Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
